# etc/services.d/ttrss-update-daemon/run
#!/usr/bin/execlineb -P
with-contenv

importas -D "" TTRSS_DB_HOST TTRSS_DB_HOST

ifelse { test -z ${TTRSS_DB_HOST} }
{
	foreground { echo "No database configured, ttrss-update-daemon stopped!" }
	foreground { touch /var/run/s6/services/ttrss-update-daemon/down }
	foreground { s6-svc -O /var/run/s6/services/ttrss-update-daemon }
	exit 1
}

s6-setuidgid nginx /usr/bin/php /var/www/update.php --daemon

