# etc/services.d/ttrss-update-daemon/run
#!/usr/bin/execlineb -P
with-contenv

s6-setuidgid nginx /usr/bin/php /var/www/update.php --daemon

