# etc/services.d/ttrss-update-daemon/run
#!/usr/bin/execlineb -P

s6-setuidgid nginx /usr/bin/php /var/www/update_daemon2.php

