# etc/s6-overlay/s6-rc.d/init-nginx/up
#!/command/execlineb -P
with-contenv

/bin/bash /usr/local/bin/nginx-conf-init.sh
