# etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
with-contenv
/usr/sbin/php-fpm81 --nodaemonize
